# supervisord.conf
[supervisord]
logfile=/var/log/supervisor.log
pidfile=/run/supervisor.pid
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
autostart=true

[include]
files=/etc/supervisor/conf.d/cert_renew.conf
